# systemd/pocket-cache-portal.service
[Unit]
Description=pocket-cache demo web portal
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/pocket-cache-firmware
ExecStart=/usr/bin/python3 portal/server.py --port 80 --bind 0.0.0.0
Environment=POCKETCACHE_CONFIG=/etc/pocket-cache/firmware.json
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
